# deploy/anika-live.service
[Unit]
Description=Anika Live Studio - 24/7 YouTube Stream Manager
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/apps/anika-live
Environment=NODE_ENV=production
Environment=PORT=3010
Environment=HOSTNAME=0.0.0.0
Environment=LIVE_DATA_ROOT=/var/www/apps/anika-live/live-data
EnvironmentFile=/var/www/apps/anika-live/.env
ExecStart=/usr/bin/node node_modules/next/dist/bin/next start -p 3010
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
